# ansible-playbooks-collection/templates/dotnet_app.service.j2
[Unit]
Description={{ app_name }} .NET Application
After=network.target

[Service]
Type=notify
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_dir }}
ExecStart=/usr/bin/dotnet {{ app_dir }}/{{ app_dll }}
Restart=always
RestartSec=10

# .NET environment
Environment="DOTNET_ENVIRONMENT={{ dotnet_environment | default('Production') }}"
Environment="ASPNETCORE_URLS=http://0.0.0.0:{{ app_port | default(5000) }}"
Environment="DOTNET_CLI_TELEMETRY_OPTOUT=1"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ app_name }}

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
