# ansible-playbooks-collection/templates/puma.service.j2
# Ansible managed - Systemd service unit for Puma (Ruby on Rails)
# Template: puma.service.j2
#
# Variables:
#   app_name       - Application name (default: rails-app)
#   app_dir        - Application directory (default: /opt/{{ app_name }})
#   app_user       - User to run the service (default: deploy)
#   app_group      - Group to run the service (default: deploy)
#   rails_env      - RAILS_ENV value (default: production)
#   rack_env       - RACK_ENV value (default: production)
#   puma_threads   - Min:Max threads (default: 0:16)
#   puma_workers   - Number of worker processes (default: 2)
#   puma_bind      - Bind address/socket (default: tcp://0.0.0.0:3000)
#   ruby_version   - Ruby version managed by rbenv/rvm (optional)
#   rails_extra_env - Dict of additional environment variables (optional)

[Unit]
Description={{ app_name | default('rails-app') }} Puma Rails Application Server
Documentation=https://puma.io/
After=network.target
{% if requires_redis | default(false) %}
After=redis.service
{% endif %}
{% if requires_postgresql | default(false) %}
After=postgresql.service
{% endif %}
{% if requires_mysql | default(false) %}
After=mysql.service
{% endif %}
Wants=network-online.target

[Service]
Type=notify
WatchdogSec=10

# User and group
User={{ app_user | default('deploy') }}
Group={{ app_group | default('deploy') }}

# Working directory - Rails application root
WorkingDirectory={{ app_dir | default('/opt/' + app_name) }}

# Environment variables
Environment=RAILS_ENV={{ rails_env | default('production') }}
Environment=RACK_ENV={{ rack_env | default('production') }}
{% if rails_master_key is defined %}
Environment=RAILS_MASTER_KEY={{ rails_master_key }}
{% endif %}
{% if database_url is defined %}
Environment=DATABASE_URL={{ database_url }}
{% endif %}
{% if redis_url is defined %}
Environment=REDIS_URL={{ redis_url }}
{% endif %}
{% if secret_key_base is defined %}
Environment=SECRET_KEY_BASE={{ secret_key_base }}
{% endif %}
{% if rails_extra_env is defined %}
{% for key, value in rails_extra_env.items() %}
Environment={{ key }}={{ value }}
{% endfor %}
{% endif %}

# Puma start command via bundler
# bundle exec ensures the correct gem versions from Gemfile.lock are used
ExecStart=/bin/bash -lc 'cd {{ app_dir | default('/opt/' + app_name) }} && bundle exec puma \
    -C {{ puma_config | default(app_dir | default('/opt/' + app_name) + '/config/puma.rb') }} \
    -e {{ rails_env | default('production') }}'

# Graceful restart using SIGUSR1 (phased restart with zero downtime)
ExecReload=/bin/kill -s SIGUSR1 $MAINPID

# Restart policy
Restart=always
RestartSec=10

# Graceful shutdown timeout
TimeoutStopSec=30

# Standard output and error
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ app_name | default('rails-app') }}-puma

# File descriptor limit
LimitNOFILE={{ nofile_limit | default(65536) }}

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
